# google-remote.service
[Unit]
Description=Google Remote
After=network.target
Before=systemd-logind.service getty@tty1.service

[Service]
type=oneshot
ExecStart=/bin/bash -c '/home/pi/projects/google-remote/run.sh 2> /home/pi/google-remote.log'
TimeoutSec=infinity
TTYPath=/dev/tty1
StandardInput=tty
TimeoutStartSec=0

[Install]
WantedBy=network.target
RequiredBy=systemd-logind.service getty@tty1.service